# configs/systemd/sift-job-worker.service
[Unit]
Description=sift-mcps Durable Job Worker
After=network.target sift-gateway.service

[Service]
Type=simple
ExecStart=${SIFT_MCPS_ROOT}/.venv/bin/sift-job-worker
WorkingDirectory=${SIFT_MCPS_ROOT}
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
EnvironmentFile=-%h/.sift/supabase.env
Environment=UV_NO_MANAGED_PYTHON=1
Environment=UV_PYTHON_DOWNLOADS=never

[Install]
WantedBy=default.target
