# configs/systemd/sift-gateway.service
[Unit]
Description=sift-mcps Gateway - MCP server runtime for digital forensics
After=network.target

[Service]
Type=simple
ExecStart=${SIFT_MCPS_ROOT}/.venv/bin/sift-gateway --config ${SIFT_CONFIG}
WorkingDirectory=${SIFT_MCPS_ROOT}
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=SIFT_EXAMINER=${SIFT_EXAMINER}

[Install]
WantedBy=default.target
