# configs/systemd/sift-job-worker.service
[Unit]
Description=sift-mcps Durable Job Worker
After=network.target sift-gateway.service

[Service]
Type=simple
# Runs as the dedicated non-admin service user (sift-service), same as the gateway.
User=${SIFT_GATEWAY_SERVICE_USER}
Group=${SIFT_GATEWAY_SERVICE_USER}
# feat/opensearch-workers: PIN this worker to the run_command lane. The shared
# job_worker_cli wires the opensearch ingest/enrich handlers whenever the add-on
# is importable, so an unrestricted worker would ALSO claim ingest/enrich jobs and
# run the FUSE mount — but this unit is calibrated for run_command (cap_linux_
# immutable only, MountFlags= empty), so the mount would fail with `sudo: unable to
# change to root gid: Operation not permitted`. The privileged ingest/enrich lane
# belongs exclusively to sift-opensearch-worker@.service (adds CAP_SYS_ADMIN and
# runs in the host mount namespace for FUSE; no MountFlags=shared). Lane isolation
# is enforced here, not just by privilege.
ExecStart=${SIFT_MCPS_ROOT}/.venv/bin/sift-job-worker --job-types run_command
WorkingDirectory=${SIFT_MCPS_ROOT}
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
# Shared Volatility3 symbol cache (mode 2775, group `sift`) — parse_memory.py /
# worker.py read SIFT_VOL_SYMBOLS and write symbols there, shared with agent_runtime.
Environment=SIFT_VOL_SYMBOLS=${SIFT_VOL_SYMBOLS}
# The worker holds the control-plane DSN and Supabase/service-role secrets it
# needs for DB job claims. These secrets MUST NOT reach a sandboxed
# run_command subprocess: that isolation is enforced in code at the execution
# boundary (sift_core.execute.runtime_acl.build_sandbox_env scrubs the env of
# both the execution worker and the forensic tool it spawns). Do not rely on
# this unit to pass secrets through to run_command — by design it cannot.
# IMPORTANT: %h would resolve to /root for a system service (NOT influenced by
# User=, see systemd.unit(5)), so we use the absolute SIFT_HOME path that
# install.sh substitutes to /var/lib/sift/.sift.
EnvironmentFile=-${SIFT_HOME}/supabase.env
EnvironmentFile=-${SIFT_HOME}/control-plane.env
# BATCH-PMI3: FK_DATA_DIR points the forensic-knowledge loader at the installed
# data dir so build_response/run_command enrichment fires under the worker user.
# Non-secret path; written by install.sh (write_fk_env). The data dir is
# read-only reference content, not a secret, so forwarding it here is safe and
# does not weaken the run_command env-scrubbing isolation described above.
EnvironmentFile=-${SIFT_HOME}/forensic-knowledge.env
Environment=UV_NO_MANAGED_PYTHON=1
Environment=UV_PYTHON_DOWNLOADS=never
# B-MVP-015 / B-MVP-004 (D3): the worker re-embeds for ingest-time RAG paths and
# memory-symbol work; point its HF cache at the service home and pin the model.
Environment=HF_HOME=${SIFT_HF_HOME}
Environment=RAG_MODEL_NAME=${SIFT_RAG_MODEL_NAME}
Environment=RAG_MODEL_REVISION=${SIFT_RAG_MODEL_REVISION}

# === BATCH-HR3 systemd hardening ============================================
# Same containment rationale as sift-gateway.service: the worker performs FUSE
# image mounts via sudo during ingest, so NoNewPrivileges / PrivateDevices /
# RestrictNamespaces / SystemCallFilter are OMITTED (would break sudo + FUSE).
# The safe, tested subset is applied below (verified: worker active, ingest mount
# + memory-symbol write to /var/cache/sift still work).
UMask=0077
ProtectSystem=strict
ReadWritePaths=/var/lib/sift /cases /var/cache/sift
ProtectHome=tmpfs
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
# Keep exactly CAP_LINUX_IMMUTABLE: the shared venv interpreter carries
# cap_linux_immutable+ep, and a file-capability binary cannot exec if that cap is
# absent from the bounding set (status=203/EXEC). All other caps are dropped.
CapabilityBoundingSet=CAP_LINUX_IMMUTABLE
# === end BATCH-HR3 hardening =================================================

[Install]
WantedBy=multi-user.target
