# configs/systemd/sift-gateway.service
[Unit]
Description=sift-mcps Gateway - MCP server runtime for digital forensics
After=network.target

[Service]
Type=simple
# Runs as the dedicated non-admin service user (sift-service). Its primary group
# is the per-user group of the same name. The aggregated stdio backend MCPs
# (opensearch-mcp/rag/opencti) are spawned as subprocesses and inherit this user.
User=${SIFT_GATEWAY_SERVICE_USER}
Group=${SIFT_GATEWAY_SERVICE_USER}
ExecStart=${SIFT_MCPS_ROOT}/.venv/bin/sift-gateway --config ${SIFT_CONFIG}
WorkingDirectory=${SIFT_MCPS_ROOT}
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=SIFT_EXAMINER=${SIFT_EXAMINER}
# Shared Volatility3 symbol cache (mode 2775, group `sift`) — written by both the
# gateway/worker (sift-service) and run_command's runtime user (agent_runtime).
Environment=SIFT_VOL_SYMBOLS=${SIFT_VOL_SYMBOLS}
# A1-BOOTSTRAP: Supabase environment variables — secrets are read from the
# process environment only, never stored in the gateway config file.
# Set SUPABASE_URL, SUPABASE_ANON_KEY, and SUPABASE_SERVICE_ROLE_KEY in
# the systemd drop-in before starting the service.
# IMPORTANT: %h is NOT usable here. For a SYSTEM service, %h resolves to /root
# (the system manager's home) and is explicitly NOT influenced by User= (see
# systemd.unit(5)). So we use the absolute SIFT_HOME path, which install.sh
# substitutes to /var/lib/sift/.sift (sift-service's passwd home + /.sift).
# These files are written by install.sh owned sift-service 0600; override
# secrets in a drop-in via:  sudo systemctl edit sift-gateway
EnvironmentFile=-${SIFT_HOME}/supabase.env
EnvironmentFile=-${SIFT_HOME}/control-plane.env
# FM-2: OpenSearch env_refs for opensearch-mcp backend process.
# OPENSEARCH_CONFIG and OPENSEARCH_HOST are resolved from the gateway process
# env and forwarded to the opensearch-mcp child via env_refs in app.mcp_backends.
# This file is written by install.sh only when SIFT_OPENSEARCH_ENABLED=true.
EnvironmentFile=-${SIFT_HOME}/opensearch.env
# BATCH-PMI3: FK_DATA_DIR points the forensic-knowledge loader at the installed
# data dir so build_response/run_command enrichment fires under the service user.
# Non-secret path; written by install.sh (write_fk_env).
EnvironmentFile=-${SIFT_HOME}/forensic-knowledge.env
# B-MVP-015 / B-MVP-004 (D3): the BGE embedding model cache lives under the
# service home (not the operator home), and the model name/revision are pinned so
# kb_search_knowledge loads reproducible weights. Non-secret identifiers/path.
Environment=HF_HOME=${SIFT_HF_HOME}
Environment=RAG_MODEL_NAME=${SIFT_RAG_MODEL_NAME}
Environment=RAG_MODEL_REVISION=${SIFT_RAG_MODEL_REVISION}

# === BATCH-HR3 systemd hardening ============================================
# Minimal, tested containment for a service that (a) spawns forensic tooling via
# `sudo -n -u agent_runtime` for run_command, and (b) mounts forensic images via
# FUSE/sudo for ingest. Those two facts forbid the strongest directives:
#   - NoNewPrivileges=yes would break `sudo` (run_command isolation) -> OMITTED.
#   - PrivateDevices/RestrictNamespaces/SystemCallFilter=@system-service can break
#     FUSE mounts and the sudo transition -> OMITTED (documented residual risk).
# What is safe and applied (verified: /health=ok, ingest mount, run_command smoke):
UMask=0077
ProtectSystem=strict
# Writable paths the service genuinely needs (everything else is read-only):
#   /var/lib/sift  — secrets/config, enrichment, hayabusa output, state
#   /cases         — evidence + agent/extraction/tmp working dirs
#   /var/cache/sift— shared volatility symbol cache
ReadWritePaths=/var/lib/sift /cases /var/cache/sift
# ProtectHome=read-only (not true): sift-service's home IS /var/lib/sift, which
# must stay writable; /home and /root are not needed and stay hidden via tmpfs.
ProtectHome=tmpfs
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
# The gateway needs no ambient capabilities of its own (sudo re-acquires what the
# child tools need through a fresh credential, outside this bounding set's reach
# for the parent). Empty bounding set drops all caps from the gateway process.
CapabilityBoundingSet=
# MemoryDenyWriteExecute is NOT set: CPython/torch JITs and some native extensions
# allocate W+X pages; enabling it crashes the embedder. Documented exclusion.
# === end BATCH-HR3 hardening =================================================

[Install]
WantedBy=multi-user.target
