# packages/forensic-knowledge/data/artifacts/windows/registry_run_keys.yaml
name: Registry Run Keys
description: "Registry Run and RunOnce keys used for automatic program execution at user logon or system startup"
platform: windows

locations:
  - path: "NTUSER.DAT\\Software\\Microsoft\\Windows\\CurrentVersion\\Run"
    os_versions: ["Windows XP+", "Windows 10", "Windows 11"]
    notes: "Per-user persistence; executes at user logon"
  - path: "NTUSER.DAT\\Software\\Microsoft\\Windows\\CurrentVersion\\RunOnce"
    os_versions: ["Windows XP+", "Windows 10", "Windows 11"]
    notes: "Per-user; executes once at next logon then deleted"
  - path: "SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run"
    os_versions: ["Windows XP+", "Windows 10", "Windows 11"]
    notes: "Machine-wide persistence (HKLM); executes for all users at logon"
  - path: "SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\RunOnce"
    os_versions: ["Windows XP+", "Windows 10", "Windows 11"]
    notes: "Machine-wide; executes once at next startup then deleted"
  - path: "SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Policies\\Explorer\\Run"
    os_versions: ["Windows XP+", "Windows 10", "Windows 11"]
    notes: "Group policy-controlled Run key; less commonly checked"
  - path: "SOFTWARE\\Wow6432Node\\Microsoft\\Windows\\CurrentVersion\\Run"
    os_versions: ["Windows 10", "Windows 11"]
    notes: "32-bit application Run key on 64-bit systems"

proves:
  - "A persistence mechanism exists or existed at the recorded path"
  - "The binary/command configured to run at startup or logon"
  - "When the registry value was last written (via key LastWriteTime)"

does_not_prove:
  - "That the configured program was actually executed (system may have been offline, user may not have logged on)"
  - "That the persistence mechanism is malicious (many legitimate applications use Run keys)"
  - "Who created the registry entry"

timestamps:
  - field: "Registry Key LastWriteTime"
    meaning: "Last time any value under the Run/RunOnce key was created, modified, or deleted"
    common_misinterpretation: "This is the key-level timestamp, not per-value. Any change to any value under the key updates this single timestamp."

common_misinterpretations:
  - claim: "Run key entries always execute"
    correction: "Run key entries execute only when the relevant logon or startup event occurs. If the user never logs on, HKCU Run entries never fire."
  - claim: "RunOnce entries are always deleted after execution"
    correction: "RunOnce entries are deleted after execution by default, but entries prefixed with '!' persist across reboots until successful execution. Entries prefixed with '*' run even in Safe Mode."
  - claim: "Checking HKLM and HKCU Run keys is sufficient"
    correction: "There are numerous additional persistence locations including RunServices, RunServicesOnce, Policies\\Explorer\\Run, and the Wow6432Node equivalents on 64-bit systems."
  - claim: "Registry key LastWriteTime tells you when a specific value was added"
    correction: "Registry LastWriteTime applies to the entire key, not individual values. Any modification to any value under the key updates the timestamp. You cannot determine which value was changed."

corroborate_with:
  for_execution: [prefetch, event_logs_security, event_logs_sysmon]
  for_presence: [mft, amcache, shimcache]
  for_timeline: [event_logs_security, usn_journal, event_logs_sysmon]

related_tools: [RECmd, autoruns]

cross_mcp_checks:
  - mcp: "windows-triage"
    tool: "check_autorun"
    when: "Persistence entry found — verify against known Windows autoruns"
  - mcp: "opencti"
    tool: "lookup_hash"
    when: "Binary hash from persistence entry — check threat intelligence"
  - mcp: "forensic-rag"
    tool: "search"
    when: "Suspicious persistence mechanism — search for related techniques and detection rules"
  - mcp: "remnux"
    tool: "analyze_file"
    when: "Unknown binary referenced by persistence entry — submit for static analysis"

references:
  - "SANS FOR500 - Windows Forensic Analysis"
  - "MITRE ATT&CK T1547.001 - Boot or Logon Autostart Execution: Registry Run Keys"

investigation_patterns:
  - name: "Scripting in autostart"
    query: "ValueData:(*powershell* OR *wscript* OR *cscript* OR *mshta* OR *cmd /c*)"
    purpose: "Run key values executing scripts — persistence via LOLBins"

  - name: "Unusual service paths"
    query: "ValueData:(*\\Temp\\* OR *\\AppData\\* OR *\\ProgramData\\* OR *\\Users\\Public\\*)"
    purpose: "Autostart entries pointing to non-standard locations — potential malware persistence"
    follow_up: "windows-triage check_file for the binary"

  - name: "WDigest credential config"
    query: "KeyPath:*WDigest* AND ValueName:UseLogonCredential"
    purpose: "WDigest enabled to store cleartext credentials in memory — credential harvesting preparation"
    note: "UseLogonCredential=1 enables cleartext credential caching. Default is disabled on Windows 8.1+."

  - name: "LSA security packages"
    query: "KeyPath:*Lsa* AND ValueName:(*Security Packages* OR *Authentication Packages*)"
    purpose: "Modified LSA security packages — SSP injection for credential interception"
